# greenhouse.cron
* * * * * root /usr/local/sbin/update_lights
